# config/default.crontab.in
# This is a sample crontab file that you can symlink to under /etc/cron.d
# in order to schedule the backup time.
#

# By default, run full xva backups every Sat. morning at 1:02AM.
2 1 * * 6 @BINDIR@/xenserver-backup >>@LOGDIR@/xenserver-backup.log 2>>@LOGDIR@/xenserver-backup.err.log
